# labgrid/deploy/exporter.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  labgrid exporter launchd service (usb-device variant).

  Uses usb-device-exporter (python -m labgrid_usb_device.run_exporter)
  instead of stock labgrid-exporter, which registers the USBDevice
  resource type for exporter.yaml.

  Before installing, update the paths below to match your environment:
    - ProgramArguments: path to your venv python and exporter.yaml
    - WorkingDirectory: directory containing exporter.yaml
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.usb-device.labgrid-exporter</string>

    <key>ProgramArguments</key>
    <array>
        <!-- Update: path to your venv python -->
        <string>/path/to/venv/bin/python</string>
        <string>-m</string>
        <string>labgrid_usb_device.run_exporter</string>
        <string>-c</string>
        <string>localhost:20408</string>
        <!-- Update: path to your exporter.yaml (copied from exporter.yaml.example) -->
        <string>/path/to/deploy/exporter.yaml</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/labgrid-exporter.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/labgrid-exporter.log</string>

    <!-- Update to the directory containing exporter.yaml -->
    <key>WorkingDirectory</key>
    <string>/path/to/deploy</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
